# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Enable and start the systemd user service for all users
        # This will be done when users log in, so we just need to ensure
        # the service is available in the user systemd directory
        if [ -x "$(command -v systemctl)" ]; then
            # Reload systemd daemon to pick up the new service file
            systemctl --global daemon-reload || true
            
            # Enable the service globally for all users
            systemctl --global enable openrgb-daemon.service || true
        fi
        ;;
    
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0